# user/.config/autostart/nm-applet.desktop

[Desktop Entry]

Type=Application
Name=network manager applet
Exec=nm-applet --indicator
