# user/.config/systemd/user/mpris-proxy.service

[Service]
ExecStart=/usr/bin/mpris-proxy

[Install]
WantedBy=default.target
